# discord-github-bot.service
[Unit]
Description=Discord GitHub Bot
After=network.target network-online.target
Wants=network-online.target

[Service]
# Set DISCORD_GITHUB_HOME to the installation directory of the bot.
# When deploying, update the path below or override the variable via an
# EnvironmentFile.
Environment=DISCORD_GITHUB_HOME=/opt/discord-github
Type=simple
User=cinder
Group=cinder
WorkingDirectory=${DISCORD_GITHUB_HOME}
Environment=PATH=${DISCORD_GITHUB_HOME}/.venv/bin
EnvironmentFile=-${DISCORD_GITHUB_HOME}/.env
ExecStart=${DISCORD_GITHUB_HOME}/.venv/bin/python run.py
Restart=always
RestartSec=10
KillMode=mixed
KillSignal=SIGINT
TimeoutStopSec=30
PrivateTmp=true
ProtectSystem=strict
ReadWritePaths=${DISCORD_GITHUB_HOME}
NoNewPrivileges=true
StandardOutput=journal
StandardError=journal
SyslogIdentifier=discord-github-bot

[Install]
WantedBy=multi-user.target
